# .devcontainer/config/supervisor/emuserver.conf
[program:emuserver]
command=/scripts/emuserver.sh
autostart=true
autorestart=true