# .devcontainer/config/supervisor/xpraserver.conf
[program:xpraserver]
command=/scripts/xpraserver.sh
autostart=true
autorestart=true